# trader-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>trader-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/binance-trader/run.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/binance-trader</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/binance-trader/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/binance-trader/bot.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
